# apps/react-drm/system/99-react-drm.rules
# Touch Bar device setup for react-drm.
# Install: sudo install -m644 system/99-react-drm.rules /etc/udev/rules.d/
#          sudo udevadm control --reload
#          sudo udevadm trigger --action=add --subsystem-match=usb --subsystem-match=backlight
#          sudo udevadm trigger --action=add --subsystem-match=misc --sysname-match=uinput
#
# Seat assignment: move the Touch Bar display and its touch input to their
# own logind seat so the main compositor never claims them.
SUBSYSTEM=="drm", KERNEL=="card*", DRIVERS=="adp|appletbdrm|t2bdrm", TAG-="master-of-seat", ENV{ID_SEAT}="seat-touchbar"
SUBSYSTEM=="input", ATTR{name}=="Apple Inc. Touch Bar Display Touchpad", ENV{ID_SEAT}="seat-touchbar"
SUBSYSTEM=="input", ATTR{name}=="MacBookPro17,1 Touch Bar", ENV{ID_SEAT}="seat-touchbar"
SUBSYSTEM=="input", ATTR{name}=="Mac14,7 Touch Bar", ENV{ID_SEAT}="seat-touchbar"

# Config switching without root: let the video group write the Touch Bar
# Display's (05ac:8302) bConfigurationValue, so the attach/detach helpers can
# flip between config 2 (react-drm display) and config 1 (firmware fn strip)
# as the session user. Re-applied automatically if the device re-enumerates.
# The devnode (GROUP/MODE) is also video-writable: the app issues a
# USBDEVFS_RESET ioctl on /dev/bus/usb/BBB/DDD before each config switch to
# wake the T2 firmware's display interface — when it has gone to sleep, every
# transfer (config writes included) fails with ETIMEDOUT until a reset.
# USB runtime PM must stay off: with the device idle in config 1,
# hid_appletb_kbd's runtime-suspend handler (appletb_kbd_set_mode) deadlocks in
# the kernel holding the device lock, which wedges any later
# bConfigurationValue write in unkillable D state until reboot. power/control=on
# alone is not enough — powertop --auto-tune flips it back to auto without a
# uevent — so autosuspend_delay_ms=-1 pins a second, independent usage-count
# reference that blocks runtime suspend even with control=auto.
#ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="05ac", ATTR{idProduct}=="8302", ATTR{bConfigurationValue}=="1", ATTR{bConfigurationValue}="0", ATTR{bConfigurationValue}="1"
ACTION=="add", SUBSYSTEM=="usb", ENV{DEVTYPE}=="usb_device", ATTR{idVendor}=="05ac", ATTR{idProduct}=="8302", GROUP="video", MODE="0660", RUN+="/bin/chown root:video /sys%p/bConfigurationValue", RUN+="/bin/chmod g+w /sys%p/bConfigurationValue", ATTR{power/control}="on", ATTR{power/autosuspend_delay_ms}="-1"

# Backlight: let the video group set Touch Bar brightness without sudo.
ACTION=="add", SUBSYSTEM=="backlight", KERNEL=="appletb_backlight|t2tb_backlight", RUN+="/bin/chown root:video /sys/class/backlight/%k/brightness", RUN+="/bin/chmod g+w /sys/class/backlight/%k/brightness"

# Key injection: let members of the input group create the virtual keyboard.
SUBSYSTEM=="misc", KERNEL=="uinput", GROUP="input", MODE="0660"
